# deploy/chatbot.service
# systemd service for ThaiLawOnline Chatbot
# Place in /etc/systemd/system/chatbot.service
#
# Commands:
#   sudo systemctl daemon-reload
#   sudo systemctl enable chatbot
#   sudo systemctl start chatbot
#   sudo systemctl status chatbot
#   journalctl -u chatbot -f

[Unit]
Description=ThaiLawOnline Chatbot API (FastAPI + LLM Council)
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/thailawonline-chatbot
EnvironmentFile=/var/www/thailawonline-chatbot/.env
ExecStart=/usr/bin/python3 -m uvicorn backend.main:app --host 127.0.0.1 --port 8001 --workers 2
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/www/thailawonline-chatbot/data

[Install]
WantedBy=multi-user.target
